# deploy/vpn-bot.nonroot.example.service
# Production sudo-helper non-root unit for vpn-bot.
# Code, deploy files, and .venv should remain root-owned and not writable by
# the vpn-bot runtime user.
# Sudo-based helpers require privilege elevation through sudo/setuid; do not
# enable the no-new-privileges hardening flag on this unit.

[Unit]
Description=VPN Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vpn-service
EnvironmentFile=/opt/vpn-service/.env
Environment=BOT_LOCK_PATH=/run/vpn-bot/vpn-bot.lock
ExecStart=/opt/vpn-service/.venv/bin/python /opt/vpn-service/main.py
Restart=on-failure
RestartSec=5
User=vpn-bot
Group=vpn-bot
UMask=0077
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
RuntimeDirectory=vpn-bot
RuntimeDirectoryMode=0700
# Sudo-launched helpers inherit this service mount namespace, so keep
# ProtectSystem=strict and expose only exact helper write targets. The account
# database files are writable only for the sudo SOCKS5 helper's useradd,
# chpasswd, passwd -l, and userdel operations.
ReadWritePaths=/opt/vpn-service/data /opt/vpn-service/logs /run/vpn-bot /usr/local/etc/xray /etc/amnezia/amneziawg /etc/mtproxy/vpnbot /etc/passwd /etc/shadow /etc/group /etc/gshadow /etc/.pwd.lock

[Install]
WantedBy=multi-user.target
