# deploy/vpn-bot.service
[Unit]
Description=VPN Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vpn-service
EnvironmentFile=/opt/vpn-service/.env
Environment=BOT_LOCK_PATH=/run/vpn-bot/vpn-bot.lock
ExecStart=/opt/vpn-service/.venv/bin/python /opt/vpn-service/main.py
Restart=on-failure
RestartSec=5
User=vpn-bot
Group=vpn-bot
UMask=0077
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
RuntimeDirectory=vpn-bot
RuntimeDirectoryMode=0700
# Sudo-launched helpers inherit this service mount namespace, so keep
# ProtectSystem=strict and expose only the exact backend trees helpers update.
# Do not enable no-new-privileges: sudo helpers require privilege elevation.
ReadWritePaths=/opt/vpn-service/data /opt/vpn-service/logs /run/vpn-bot /usr/local/etc/xray /etc/amnezia/amneziawg /etc/mtproxy/vpnbot

[Install]
WantedBy=multi-user.target
