# deploy/vpnbot-warp-split.service
[Unit]
Description=vpnbot WARP selective-split (listed prefixes via WARP, rest direct)
After=warp-routes.service network-online.target
Wants=network-online.target
PartOf=warp-routes.service
ConditionPathExists=/etc/vpnbot/warp-split.list

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/vpnbot-warp-split apply
ExecStop=/usr/local/sbin/vpnbot-warp-split revert

[Install]
WantedBy=multi-user.target
