# deploy/warp-routes.service
# WARP policy-routing rules — applied at boot and torn down on stop.
#
# Diverts only the AmneziaWG client subnet (10.0.0.0/24) out through the
# out-warp tunnel while the host itself stays on the direct path. The tunnel is
# brought up by awg-quick@out-warp with `Table = auto`, which creates the dynamic
# routing table and the host-bypass rules; this unit's ExecStart reads that table
# at runtime, STRIPS the host-bypass (so SSH/the bot keep egressing directly) and
# installs the narrow client rule, the anti-loop endpoint pin, the out-warp
# MASQUERADE and the FORWARD accepts above UFW. The helper is idempotent on add
# and safe on del, so this unit can be (re)started freely.
#
# The helper runs directly as root here (no sudo). systemd owns both the
# interface (awg-quick@out-warp) and these routes; the bot's health monitor is a
# pure observer and never adds/removes them. Binding to awg-quick@out-warp
# (Requires=/After=/PartOf=) guarantees the routes are installed only after the
# interface is up and re-applied whenever it restarts, so the two owners can
# never drift apart. awg-quick@awg0 (the client interface the FORWARD rules name)
# is ordered first as well.
#
# Install (interfaces first, then the routes that ride on them):
#   install -o root -g root -m 0755 scripts/vpnbot-warp-routes /usr/local/sbin/vpnbot-warp-routes
#   install -o root -g root -m 0644 deploy/warp-routes.service /etc/systemd/system/warp-routes.service
#   systemctl daemon-reload
#   systemctl enable --now awg-quick@out-warp
#   systemctl enable --now warp-routes.service
[Unit]
Description=WARP policy routing rules
After=network-online.target awg-quick@out-warp.service awg-quick@awg0.service
Wants=network-online.target awg-quick@awg0.service
Requires=awg-quick@out-warp.service
PartOf=awg-quick@out-warp.service
Before=vpn-bot.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/vpnbot-warp-routes add out-warp
ExecStop=/usr/local/sbin/vpnbot-warp-routes del out-warp

[Install]
WantedBy=multi-user.target
